# arc-cms.service
[Unit]
Description=ARC CMS Django Application with uWSGI
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=notify
User=ubuntu
Group=ubuntu
# WorkingDirectory will be set by the deployment script
Environment=DJANGO_SETTINGS_MODULE=cms_core.settings
# PYTHONPATH will be set dynamically
ExecStart=/bin/bash -c 'cd "$ARC_CMS_DIR" && "$ARC_CMS_DIR/venv/bin/uwsgi" --ini uwsgi.ini'
ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s QUIT $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
# ReadWritePaths will be set dynamically

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
