# deployment/systemd/arc-cms.service
[Unit]
Description=ARC CMS Django Application with uWSGI
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=notify
User=${DEPLOY_USER}
Group=${DEPLOY_GROUP}
WorkingDirectory=${BACKEND_DIR}
Environment=DJANGO_SETTINGS_MODULE=cms_core.settings
Environment=PYTHONPATH=${BACKEND_DIR}
Environment=ARC_CMS_DIR=${BACKEND_DIR}
ExecStart=${BACKEND_DIR}/venv/bin/uwsgi --ini uwsgi.ini
ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s QUIT $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=${BACKEND_DIR}
ReadWritePaths=${LOG_DIR}
ReadWritePaths=/var/run
ReadWritePaths=/tmp

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
